# config/backup/backup.timer
[Unit]
Description=Backup timer

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
