# config/nextcloud/preview-pregenerate.timer
[Unit]
Description = Nextcloud preview pregenerate

[Timer]
# Daily
OnCalendar = *-*-* 02:00:00

[Install]
WantedBy = timers.target
