# config/updates/update-security-packages.timer
[Unit]
Description = Update security packages

[Timer]
# Daily
OnCalendar = *-*-* 04:00:00
RandomizedDelaySec = 5m

[Install]
WantedBy = timers.target
